# Cloud/Docker/Dockerfile/Togetherjs/supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:wsserver]  
command=node /togetherjs/hub/server.js 
  
[program:site]  
command=node /togetherjs/devserver.js
